# .config/systemd/user/dotfiles-update.timer
[Unit]
Description=update dotfiles daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
