# .config/systemd/user/dotfiles-update.service
[Service]
Description=update dotfiles
WorkingDirectory=%h
ExecStart=make update

[Install]
WantedBy=default.target

