# .config/systemd/user/udiskie.service
[Service]
Description=automounter
ExecStart=udiskie -an
Restart=always

[Install]
WantedBy=default.target

