# .config/systemd/user/udiskie.service
[Service]
Description=automounter
ExecStart=udiskie -an -F --no-terminal --appindicator
Restart=always

[Install]
WantedBy=default.target

